# server/deploy/outfit-advisor.service
# Outfit Advisor — DGX /advice server (systemd USER service, no sudo needed)
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp outfit-advisor.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now outfit-advisor
#   loginctl enable-linger korety        # keep it running after logout / across reboots
#
# Verify:  curl http://100.112.171.54:8787/health   ->  {"ok":true,"vllm":...}
# Logs:    journalctl --user -u outfit-advisor -f

[Unit]
Description=Outfit Advisor /advice server (FastAPI + vLLM)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/korety/Project/outfit-advisor/server
# SECURITY: bind the Tailscale IP, not 0.0.0.0 — the LAN must never see this port.
# If tailscale0 isn't up yet at boot, the bind fails and Restart retries until it is.
ExecStart=/home/korety/miniconda3/bin/python3 -m uvicorn app:app --host 100.112.171.54 --port 8787 --no-access-log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
